# systemd/clamscan.timer
[Unit]
Description=Daily scan of clamscan

[Timer]
OnCalendar=*-*-* 23:59:59
Persistent=true

[Install]
WantedBy=timers.target
