# systemd/mngip.timer
[Unit]
Description=Hourly scan of management global ip

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
